# project/scripts/package/DEBIAN/postrm
#!/bin/bash

update-rc.d topicator remove

rm -fr /etc/topicator
rm -fr /opt/topicator

