# project/scripts/package/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/topicator

update-rc.d topicator defaults

service topicator start